# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw AI Gateway
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service
StartLimitBurst=5
StartLimitIntervalSec=600

[Service]
Type=simple
User=henry
WorkingDirectory=/opt/openclaw
EnvironmentFile=-/opt/openclaw/.env
ExecStart=/usr/bin/openclaw gateway start
ExecStop=/usr/bin/openclaw gateway stop
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/opt/openclaw /var/log/openclaw

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

[Install]
WantedBy=multi-user.target
